# launchd/com.btcbot.phase0_settle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.btcbot.phase0_settle</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/brooksmoore/Desktop/umbrella/.venv/bin/python</string>
		<string>/Users/brooksmoore/Desktop/btc-bot/Short-term Bitcoin Contract Trader Bot/deploy/run_phase0_settlements.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/brooksmoore/Desktop/btc-bot/Short-term Bitcoin Contract Trader Bot</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>/Users/brooksmoore/Desktop/btc-bot/Short-term Bitcoin Contract Trader Bot:/Users/brooksmoore/Desktop/umbrella</string>
		<key>KALSHI_ENV</key>
		<string>prod</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<!-- Every 5 minutes -->
	<key>StartInterval</key>
	<integer>300</integer>
	<key>StandardOutPath</key>
	<string>/tmp/btc_phase0_settle_launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/btc_phase0_settle_launchd.log</string>
</dict>
</plist>
